# at-boot.target
[Unit]
Description=At-Boot Mode
Requires=basic.target at-boot.service
AllowIsolate=yes
